# roles/logstash-exporter/templates/logstash_exporter.service.j2
[Unit]
Description=Logstash Exporter
After=network.target

[Service]
ExecStart=/usr/local/bin/logstash-exporter {{ logstash_exporter_options }}
Restart=always

[Install]
WantedBy=multi-user.target
